# tw-repo/cookbooks/thoughtworks/templates/default/jar_service.erb
#!/bin/sh

if [ true != "$INIT_D_SCRIPT_SOURCED" ] ; then
  set "$0" "$@"; INIT_D_SCRIPT_SOURCED=true . /lib/init/init-d-script
fi

### BEGIN INIT INFO
# Provides: <%= @appName %>
# Description: <%= @appName %> service
# Required-Start: $network $local_fs $syslog
# Default-Start: 2 3 4 5
# Default-Stop: 0 1 6 
### END INIT INFO

# Author: Matt Hines <matt@matthin.es>

DESC="Starts and stops <%= @appName %> service"

NAME=<%= @appName %>
JAR_HOME=<%= @appHome -%>/$NAME
PATH_TO_JAR=$JAR_HOME/$NAME.jar
export APP_PORT=8080
<% if @appName == "front-end" %>
export STATIC_URL=localhost
export QUOTE_SERVICE_URL=http://<%= node['thoughtworks']['quotes_ip'] %>:8080
export NEWSFEED_SERVICE_URL=http://<%= node['thoughtworks']['newsfeed_ip'] %>:8080
export NEWSFEED_SERVICE_TOKEN="<%= node['thoughtworks']['newsfeed_api_token'] %>"
<% end %>

PIDFILE=/var/run/$NAME.pid
LOGFILE=/var/log/$NAME.log

RETVAL=0

start() {
    if [ -f $PIDFILE ] && kill -0 $(cat $PIDFILE); then
        echo "$NAME service already running ..."
    else
        echo "Starting $NAME service ..."
        cd $JAR_HOME && java -jar $PATH_TO_JAR >> $LOGFILE 2>&1 &
        echo $! > $PIDFILE
        echo "$NAME service started ..."
    fi
}

stop() {
    if [ -f $PIDFILE ]; then
      PID=$(cat $PIDFILE);
      echo "Stopping $NAME service ($PID) ..."
      kill $PID
      echo "$NAME service stopped ..."
      rm $PIDFILE
    else
       echo "$NAME service is not running ..."
    fi
}

restart() {
    stop
    sleep 3
    start
}

status() {
    if [ -f $PIDFILE ]; then
        echo "$NAME service running with pid: $(cat $PIDFILE); ..."
        return 0
    else
        echo "$NAME service not running ..."
        return 3
    fi
}

case "$1" in
    start)
        start
        ;;
    stop)
        stop
        ;;
    restart)
        restart
        ;;
    status)
        status
        ;;
    *)
        echo "Usage $0 {start|stop|restart|status}"
        RETVAL=1
esac
